# deploy/output/amd-vrt-2025-08-05-14-36-53/DEBIAN/postrm
#!/bin/bash
set -e

# Remove module loading configuration
if [ -f "/etc/modules-load.d/amd-vrt.conf" ]; then
    rm -f /etc/modules-load.d/amd-vrt.conf
fi

# Remove udev rules
if [ -f "/etc/udev/rules.d/99-amd-vrt-permissions.rules" ]; then
    rm -f /etc/udev/rules.d/99-amd-vrt-permissions.rules
    udevadm control --reload-rules
fi

# Remove systemd service
if [ -f "/etc/systemd/system/vrt-devices.service" ]; then
    systemctl disable vrt-devices.service || true
    systemctl stop vrt-devices.service || true
    rm -f /etc/systemd/system/vrt-devices.service
    systemctl daemon-reload
fi

# Remove device setup script
if [ -f "/usr/local/bin/vrt-setup-devices.sh" ]; then
    rm -f /usr/local/bin/vrt-setup-devices.sh
fi

# Remove ldconfig configuration
if [ -f "/etc/ld.so.conf.d/amd-vrt.conf" ]; then
    rm -f /etc/ld.so.conf.d/amd-vrt.conf
    ldconfig
fi

exit 0
